# package/ai-monitor/files/ai-monitor.init
#!/bin/sh /etc/rc.common

START=95
STOP=10
USE_PROCD=1
PROG=/usr/bin/ai-monitor.sh
CONF=/etc/ai-monitor.conf
LIB_DIR=/usr/lib/ai-monitor

uci_to_conf() {
    mkdir -p /etc
    cat > "$CONF" <<EOF
# Generated by UCI - do not edit manually
CHECK_INTERVAL=$(uci -q get ai-monitor.general.check_interval || echo 300)
AI_API_URL=$(uci -q get ai-monitor.ai.api_url || echo "https://api.deepseek.com/chat/completions")
AI_API_KEY=$(uci -q get ai-monitor.ai.api_key || echo "")
AI_MODEL=$(uci -q get ai-monitor.ai.model || echo "deepseek-chat")
PUSH_TYPE=$(uci -q get ai-monitor.push.type || echo "none")
PUSH_TOKEN=$(uci -q get ai-monitor.push.token || echo "")
TELEGRAM_CHAT_ID=$(uci -q get ai-monitor.push.telegram_chat_id || echo "")
CPU_THRESHOLD=$(uci -q get ai-monitor.thresholds.cpu || echo 80)
MEM_THRESHOLD=$(uci -q get ai-monitor.thresholds.memory || echo 85)
DISK_THRESHOLD=$(uci -q get ai-monitor.thresholds.disk || echo 90)
TEMP_THRESHOLD=$(uci -q get ai-monitor.thresholds.temperature || echo 75)
REPORT_TO_PUSH=$(uci -q get ai-monitor.report.push_enabled || echo 1)
AUTO_OPTIMIZE=$(uci -q get ai-monitor.optimize.auto_enabled || echo 0)
OPTIMIZE_INTERVAL=$(uci -q get ai-monitor.optimize.interval || echo 3600)
REPORT_DAILY_HOUR=$(uci -q get ai-monitor.report.daily_hour || echo 8)
REPORT_WEEKLY_DAY=$(uci -q get ai-monitor.report.weekly_day || echo 1)
REPORT_WEEKLY_HOUR=$(uci -q get ai-monitor.report.weekly_hour || echo 8)
EOF
}

start_service() {
    [ -x "$PROG" ] || return 1
    uci_to_conf
    procd_open_instance
    procd_set_param command "$PROG"
    procd_set_param respawn 3600 5 86400
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

reload_service() {
    uci_to_conf
    restart
}

service_triggers() {
    procd_add_reload_trigger "ai-monitor"
}

stop_service() {
    killall ai-monitor.sh 2>/dev/null
}
