# package/luci-app-ai-monitor/root/etc/uci-defaults/90-ai-monitor
#!/bin/sh
uci -q batch <<EOF
set ai-monitor.general=general
set ai-monitor.general.check_interval=300
set ai-monitor.ai=ai
set ai-monitor.ai.api_url=https://api.deepseek.com/chat/completions
set ai-monitor.ai.model=deepseek-chat
set ai-monitor.push=push
set ai-monitor.push.type=serverchan
set ai-monitor.thresholds=thresholds
set ai-monitor.thresholds.cpu=80
set ai-monitor.thresholds.memory=85
set ai-monitor.thresholds.disk=90
set ai-monitor.thresholds.temperature=75
set ai-monitor.report=report
set ai-monitor.report.push_enabled=1
set ai-monitor.report.daily_hour=8
set ai-monitor.report.weekly_day=1
set ai-monitor.report.weekly_hour=8
set ai-monitor.optimize=optimize
set ai-monitor.optimize.auto_enabled=0
set ai-monitor.optimize.interval=3600
commit ai-monitor
EOF
exit 0
